# provisioning/setup/mailhog.service
[Unit]
Description=MailHog for local
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/root/go/bin/MailHog
[Install]
WantedBy=multi-user.target